# config/systemd/user/gdisk-sync.timer
[Unit]
Description=Synchronize gdisk files every hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target

